# raspberry_pi/sounds-sync.timer
[Unit]
Description=Run Proton Drive sound sync every 5 minutes

[Timer]
OnBootSec=30s
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
